# leash.service
[Unit]
Description=Leash - AI Agent Visibility
Documentation=https://leash.meridianhouse.tech
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/leash start
ExecStop=/usr/local/bin/leash stop
Restart=always
RestartSec=3

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/leash %h/.local/state/leash %h/.config/leash
PrivateTmp=true
NoNewPrivileges=false
CapabilityBoundingSet=CAP_SYS_PTRACE CAP_KILL CAP_NET_ADMIN
AmbientCapabilities=CAP_SYS_PTRACE

# Resource limits
MemoryMax=256M
CPUQuota=25%

[Install]
WantedBy=multi-user.target
